# deployment/race-cycling-api.service
[Unit]
Description=Race Cycling History API Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/race-cycling-app
Environment=PATH=/var/www/race-cycling-app/venv/bin
Environment=DB_PATH=/var/www/race-cycling-app/backend/database/cycling_data.db
Environment=AUTH_DB_PATH=/var/www/race-cycling-app/backend/database/auth.db
Environment=FLASK_ENV=production
ExecStart=/var/www/race-cycling-app/venv/bin/python backend/api/server.py
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target